# x11-drivers/displaylink-driver/files/rc-displaylink-1.3
#!/sbin/openrc-run

depend() {
	after display-manager
}

start_pre() {
	modprobe -q evdi
}

start() {
	start-stop-daemon --start --background --chdir /opt/displaylink --make-pidfile --pidfile /run/displaylink.pid --exec /opt/displaylink/DisplayLinkManager
	eend $?
}

stop() {
	start-stop-daemon --stop --pidfile /run/displaylink.pid
	eend $?
}
